# pe-format.init.in
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

extra_commands=reload
description='Associate PE binaries with pe-exec using binfmt_misc'
description_reload='Regenerate the configuration file for pe-exec'

depend() {
	after procfs
	use localmount
}

reload() {
	ebegin "Generating configuration for pe-exec"
	export UNKNOWN MSDOS WIN32 CLR WIN64
	@sbindir@/pe-format2-conf.sh
	eend ${?} || return ${?}
}

start() {
	reload || return ${?}

	ebegin "Registering PE binaries with pe-exec wrapper"

	if [ ! -d /proc/sys/fs/binfmt_misc ]; then
		eend 1
		eerror 'You need to enable support for "misc binaries" in your kernel!'
		return 1
	fi

	if [ -f /proc/sys/fs/binfmt_misc/PE ]; then
		eindent
		ebegin "Unregistering another handler for PE files."
		echo '-1' > /proc/sys/fs/binfmt_misc/PE
		eend ${?}
		ebegin "Registering the new handler"
		eoutdent
	elif [ ! -f /proc/sys/fs/binfmt_misc/register ]; then
		local ret

		mount -t binfmt_misc binfmt_misc /proc/sys/fs/binfmt_misc
		ret=${?}
		if [ ${ret} -ne 0 ]; then
			eend ${ret}
			return ${ret}
		fi
	fi

	echo ':PE:M::MZ::/usr/bin/pe-exec:' > /proc/sys/fs/binfmt_misc/register
	eend ${?}
}

stop() {
	ebegin "Unregistering PE binaries"

	if [ -f /proc/sys/fs/binfmt_misc/PE ]; then
		echo '-1' > /proc/sys/fs/binfmt_misc/PE
	fi

	eend ${?}
}
